# app-laptop/touchbard/files/touchbard
#!/sbin/openrc-run

command="/usr/bin/touchbard"
command_background=true
pidfile="/run/touchbard/touchbard.pid"
