# app-laptop/touchbard/files/touchbard
#!/sbin/openrc-run

command="/usr/bin/touchbard"
pidfile="/run/touchbard.pid"

start() {
	ebegin "Starting t2fand"

	start-stop-daemon --start --exec ${command} \
		--pidfile ${pidfile} --background --
	eend $?
}
