# app-laptop/t2fand/files/t2fand
#!/sbin/openrc-run

command="/usr/bin/t2fand"
pidfile="/run/t2fand.pid"

start() {
	ebegin "Starting t2fand"

	start-stop-daemon --start --exec ${command} \
		--pidfile ${pidfile} --background --
	eend $?
}
