# app-laptop/t2fand/files/t2fand
#!/sbin/openrc-run

command="/usr/bin/t2fand"
command_background=true
pidfile="/run/t2fand.pid"
